# docker/php8-fpm/supervisor/cron.conf
[program:cron]
command=cron -f
umprocs=1
autostart=true
autorestart=true
startretries=10
user = root
group = root
startsecs = 0
stdout_logfile=/var/www/runtime/cron.log